# Deploy/boxaroo-long-categories.service
[Unit]
Description=Run Boxaroo LONG categories sequentially
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/odyssey/boxaroo
ExecStart=/bin/bash -lc 'set -euo pipefail; mapfile -t categories < <(/home/odyssey/boxaroo/.venv/bin/python -c "import json, pathlib; data = json.loads(pathlib.Path(\"Data/category_lists/woolworths-category-lists.json\").read_text()); print(\"\\n\".join(data[\"long\"]))"); for category in "${categories[@]}"; do echo "[boxaroo] Running category: ${category}"; /home/odyssey/boxaroo/.venv/bin/python __main__.py --category "${category}" --headless; done'
User=odyssey
