# Deploy/boxaroo-long-categories.example.service
[Unit]
Description=Run Boxaroo LONG categories sequentially
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/<your-user>/boxaroo
ExecStart=/bin/bash /home/<your-user>/boxaroo/Scripts/run_long_categories.sh
User=<your-user>
